# scripts/systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway Service
Documentation=https://docs.openclaw.ai/gateway
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Environment="NODE_ENV=production"
Environment="PATH=/home/codespace/.local/bin:/home/codespace/.npm-global/bin:/home/codespace/bin:/home/codespace/.nvm/current/bin:/home/codespace/nvm/current/bin:/home/codespace/.fnm/current/bin:/home/codespace/.volta/bin:/home/codespace/.asdf/shims:/home/codespace/.local/share/pnpm:/home/codespace/.bun/bin:/usr/local/share/nvm/current/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
WorkingDirectory=/workspaces/openclaw
ExecStart=/home/codespace/nvm/current/bin/node /workspaces/openclaw/dist/entry.js gateway --port 18789
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-gateway

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/codespace/.openclaw

[Install]
WantedBy=default.target
