# scripts/systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway Service
Documentation=https://docs.openclaw.ai/gateway
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Environment="NODE_ENV=production"
# NOTE: This unit file is intended for GitHub Codespaces; update the WorkingDirectory if used elsewhere.
WorkingDirectory=/workspaces/openclaw
ExecStart=/bin/bash -lc 'exec openclaw gateway run --bind loopback --port 18789 --force'
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-gateway

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/codespace/.openclaw

[Install]
WantedBy=default.target
